# config/simple_scheduler.service
[Unit]
Description=Simple program scheduler

[Service]
User=root
ExecStart=simple_scheduler

[Install]
WantedBy=multi-user.target